# rootfs/etc/s6-overlay/s6-rc.d/update-text/run
#!/command/with-contenv sh
exec /usr/bin/python3 /update-text.py