# rootfs/etc/s6-overlay/s6-rc.d/now-playing-yt/run
#!/command/with-contenv sh
exec /usr/bin/python3 /now-playing-yt.py